# dot_config/systemd/user/daily-memo-sync.timer
[Unit]
Description=Daily Memo Git Sync Timer

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
