# dot_config/systemd/user/daily-memo-sync.service
[Unit]
Description=Daily Memo Git Sync

[Service]
Type=oneshot
ExecStart=/home/emacs/.local/bin/daily-memo-sync
