# dot_config/systemd/user/gcalcli-remind.timer
[Unit]
Description=Google Calendar Reminder Timer

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min

[Install]
WantedBy=timers.target
