# dot_config/systemd/user/mochimochi-publish.timer
[Unit]
Description=もちもちAIトピックス daily publish

[Timer]
OnCalendar=*-*-* 07:00:00
OnCalendar=*-*-* 12:00:00
OnCalendar=*-*-* 18:00:00
Persistent=true

[Install]
WantedBy=timers.target
