# dot_config/systemd/user/mochimochi-engineer.timer
[Unit]
Description=もちもちエンジニアトピックス publish

[Timer]
OnCalendar=Tue..Sun *-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
