# dot_config/systemd/user/mochimochi-publish.timer
[Unit]
Description=もちもちAIトピックス publish

[Timer]
OnCalendar=Mon *-*-* 07:00:00
Persistent=true

[Install]
WantedBy=timers.target
